# finflow-queue.service
[Unit]
Description=FinFlow Edge Queue Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/opc/finflow-queue
EnvironmentFile=/home/opc/finflow-queue/finflow-queue.env
ExecStart=/home/opc/finflow-queue/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=opc

[Install]
WantedBy=multi-user.target
